# bedrock-timeline.service
[Unit]
Description=Artisan's Timeline - Minecraft Bedrock Player Activity Monitor
After=network.target

[Service]
Type=simple
User=bedrock
Group=bedrock
WorkingDirectory=/opt/bedrock-timeline
ExecStart=/opt/bedrock-timeline/bedrock-timeline
Restart=always
RestartSec=5

EnvironmentFile=/opt/bedrock-timeline/.env

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target